# systemd/profile.d/pufferpanel.sh
export PUFFER_CONFIG=/etc/pufferpanel/config.json